# ai-token-usage/com.user.openrouter-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.openrouter-scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rede/git/sh/ai-token-usage/scrape_openrouter.sh</string>
    </array>

    <!-- 每天10:00运行（监控前一天数据） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 开机后延迟10分钟自动运行一次 -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rede/git/sh/ai-token-usage/logs/openrouter_launchd_output.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rede/git/sh/ai-token-usage/logs/openrouter_launchd_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
