# natural-gas/com.user.natural-gas-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.natural-gas-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rede/git/sh/natural-gas/get_natural_gas.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rede/git/sh/natural-gas/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rede/git/sh/natural-gas/logs/launchd.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/rede/git/sh/natural-gas</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
